# demo/clock/epd_clock.service
[Unit]
Description= epaper display clock

[Service]
Type=oneshot
WorkingDirectory=
User=
ExecStart=

